# assets/rpi0wcar.service
[Unit]
Description=RaspberryPi Zero W Controlled RC car

[Service]
ExecStart=/usr/local/bin/rpi0wcar

[Install]
WantedBy=multi-user.target